# ops/cms-auth-proxy/systemd-example.service
[Unit]
Description=Decap CMS GitHub OAuth Proxy
After=network.target

[Service]
Type=simple
WorkingDirectory="/home/tomakin/Desktop/Portfolio Website/app/ops/cms-auth-proxy"
EnvironmentFile="/home/tomakin/Desktop/Portfolio Website/app/ops/cms-auth-proxy/.env"
ExecStart=/usr/bin/node "/home/tomakin/Desktop/Portfolio Website/app/ops/cms-auth-proxy/server.mjs"
Restart=always
RestartSec=2
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
